# qsdk/package/network/config/netifd/files/etc/hotplug.d/iface/98-wdsDown
[ "$ACTION" = ifdown ] || exit 0
[ "$INTERFACE" = wwan ] || exit 0
echo '##########################wds down!############################' > /dev/console
wdsctrl