# qsdk/package/network/config/netifd/files/etc/hotplug.d/iface/99-wdsWhenEthDown
#!/bin/sh
[ "$ACTION" = ifdown ] || exit 0
[ "$INTERFACE" = wan ] || exit 0
ifup wwan