# qsdk/package/base-files/files/etc/init.d/smd
#!/bin/sh /etc/rc.common
# Copyright (C) 2014 OpenWrt.org

START=20

start()
{
    ulimit -c unlimited
    smd

}

stop()
{
    killall smd
}
